# flipper/challenge/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/bin/gunicorn flipper:app -b localhost:5000 -w 8
directory=/srv
user=flipper

[program:nginx]
command=/usr/sbin/nginx
