# randhum/challenge/crontab
*/10 * * * * root dd if=/dev/urandom of=/tmp/rngseed bs=4096 count=1
